# radicale.cron
*/3 * * * *  exec /home/ubuntu/bin/onlyone /home/ubuntu/bin/radicale earl
*/3 * * * *  exec /home/ubuntu/bin/onlyone /home/ubuntu/bin/sshtunnel-52321-web.timberdragon.com
