# radicale.cron
*/3 * * * *  exec "$HOME/bin/onlyone" "$HOME/bin/radicale" earl
*/3 * * * *  exec "$HOME/bin/onlyone" "$HOME/bin/sshtunnel-52321-web.timberdragon.com"
